# batt_low.conf
[program:batt_low]
command="/usr/local/bin/batt_low.py"
