# gpio_shutdown.conf
[program:gpio_shutdown]
command="/usr/local/bin/gpio_shutdown.py"
