# gpio_getty.conf
[program:gpio_getty]
command="/usr/local/bin/gpio_getty.py"
autorestart=unexpected
